# Software/services/vision.service
#-----------------------------------------------------------------
# Vision Processing Service
#-----------------------------------------------------------------
[Unit]
Description=Vision Processing
After=network.target

[Service]
Type=exec
User=root
ExecStart=/bin/bash /home/pi/prog/software/scripts/launchVizPipeRobot.sh
KillMode=mixed
TimeoutStopSec=10

# Resource Limits
CPUQuota=60%
MemoryMax=256M
MemorySwapMax=64M
IOWeight=150 
Nice=-20 

[Install]
WantedBy=multi-user.target